# systemd/spotify-kiosk.service.in
[Unit]
Description=Spotify Now Playing Kiosk
Wants=spotify-now-playing.service
After=systemd-user-sessions.service spotify-now-playing.service

[Service]
User=@KIOSK_USER@
PAMName=login
TTYPath=/dev/tty7
StandardInput=tty
StandardOutput=journal
StandardError=journal
Environment=KIOSK_URL=@KIOSK_URL@
Environment=XDG_RUNTIME_DIR=/run/user/%U
ExecStartPre=/bin/mkdir -p /run/user/%U
ExecStartPre=/bin/chown %U:%U /run/user/%U
ExecStartPre=/bin/chmod 700 /run/user/%U
ExecStart=/usr/bin/dbus-run-session /usr/bin/labwc
Restart=on-failure
RestartSec=3

[Install]
WantedBy=graphical.target
